# packaging/forza-dualsense-backend.service
[Unit]
Description=Forza DualSense SteamOS backend
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/forza-dualsense/backend
ExecStart=/usr/bin/python3 %h/.local/share/forza-dualsense/backend/run_backend.py run
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
